# mysqlconf/mariadb.socket
[Unit]
Description=MariaDB Socket for local connections

[Socket]
ListenStream=/opt/server/mysql/tmp/mysqld.sock

[Install]
WantedBy=sockets.target
